# habit_tracking.cron
59 23 * * sun    {PUT HABIT TRACKING DIRECTORY HERE}/habits_cron.sh
